# launchd/com.mathur-prerit.agentisland.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.mathur-prerit.agentisland</string>

  <!-- Replace the path below with your built agentislandd binary, e.g.
       /opt/homebrew/bin/agentislandd  (Homebrew)  or
       /Users/<you>/projects/prerit/agent-island/.build/release/agentislandd  (local build).
       Then: cp this file to ~/Library/LaunchAgents/ and `launchctl load` it. -->
  <key>ProgramArguments</key>
  <array>
    <string>__REPLACE_WITH_ABSOLUTE_PATH_TO__/agentislandd</string>
  </array>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
